# cpu-thermal-guard.service
[Unit]
Description=CPU Thermal Guard (sustain-gated NBFC fan pulser)
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/cpu-thermal-guard --config /etc/cpu-thermal-guard.conf
Restart=on-failure
RestartSec=5
Nice=10
IOSchedulingClass=idle
# Hardening
ProtectHome=yes
ProtectSystem=strict
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
